# systemd/llama-qwen3-vl.service
[Unit]
Description=llama.cpp server for Qwen3-VL-8B (Vision - Computer Use)
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%HOME%/llama.cpp
Environment="CUDA_VISIBLE_DEVICES=0"
ExecStart=%HOME%/llama.cpp/build/bin/llama-server \
    --model %HOME%/models/qwen3-vl-8b/Qwen_Qwen3-VL-8B-Instruct-Q4_K_M.gguf \
    --mmproj %HOME%/models/qwen3-vl-8b/mmproj-Qwen_Qwen3-VL-8B-Instruct-f16.gguf \
    --ctx-size 4096 \
    --port 8084 \
    --n-gpu-layers 99 \
    --threads %NPROC% \
    --host 127.0.0.1
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
