# systemd/llama-qwen3-embedding.service
[Unit]
Description=llama.cpp server for Qwen3 Embedding 8B (embeddings, port 8087)
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%HOME%/llama.cpp
Environment="CUDA_VISIBLE_DEVICES=0"
LimitMEMLOCK=infinity
ExecStart=%HOME%/llama.cpp/build/bin/llama-server \
    --model /mnt/nvme/models/qwen3-embedding-8b/Qwen3-Embedding-8B-Q4_K_M.gguf \
    --ctx-size 8192 \
    --port 8087 \
    --n-gpu-layers 99 \
    --threads %NPROC% \
    --host 127.0.0.1 \
    --embeddings \
    --pooling last \
    --mlock \
    --alias qwen3-embedding-8b
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
