# systemd/llama-minicpm-v.service
[Unit]
Description=llama.cpp server for MiniCPM-V 2.5 (Vision)
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%HOME%/llama.cpp
Environment="CUDA_VISIBLE_DEVICES=0"
ExecStart=%HOME%/llama.cpp/build/bin/llama-server \
    --model %HOME%/models/minicpm-v/ggml-model-Q4_K_M.gguf \
    --mmproj %HOME%/models/minicpm-v/mmproj-model-f16.gguf \
    --ctx-size 4096 \
    --port 8083 \
    --n-gpu-layers 99 \
    --threads %NPROC% \
    --host 127.0.0.1
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
